# codedeploy/restartwebapp.service
[Unit]
Description=Restart the REST Service

[Service]
User=ubuntu

WorkingDirectory=/var
ExecStart=/var/applicationstart.sh
Restart=always

[Install]
WantedBy=multi-user.target